# dot_config/systemd/user/enwiro-daemon.service.d/private_override.conf.tmpl
[Service]
Environment=GH_TOKEN={{ output "gh" "auth" "token" | trim }}
